# openrc/backlight
#!/sbin/openrc-run
# created by jacopo@autistici.org, with patch from Raphael Costa Raphael
# with this software you have no kind of support.
#

start() {
  ebegin 'Restoring the screen brightness'

    for card in `find /sys/class/backlight/ -type l`; do
      if [[ -r "/var/cache/backlight/`basename $card`-brightness-old" ]]; then
          cp "/var/cache/backlight/`basename $card`-brightness-old" "/sys/class/backlight/`basename $card`/brightness"
      fi
    done

  eend $?
}

stop() {
  ebegin 'Saving the screen brightness'
  if ! [[ -d '/var/cache/backlight/' ]]; then
    mkdir '/var/cache/backlight'
  fi
  if ! [[ -w '/var/cache/backlight/' ]]; then
    chmod 755 '/var/cache/backlight/'
  fi

  for card in `find /sys/class/backlight/ -type l`; do
    cp "/sys/class/backlight/`basename $card`/brightness" "/var/cache/backlight/`basename $card`-brightness-old"
  done
  eend $?
}
